# daemon/claude-channel-agent.service
# Headless switchboard CHANNEL responder — systemd USER unit.
# A persistent, full-context `claude --channels` session (vs. the cold
# claude-code-agent daemon's per-message `claude --print`). Install:
#   cp to ~/.config/systemd/user/claude-channel-agent.service
#   systemctl --user daemon-reload && systemctl --user enable --now claude-channel-agent
#
# PREREQUISITE: Claude CLI authenticated (claude.ai OAuth, full scope) on this host,
# workspace trust pre-accepted for CHANNEL_WORKDIR, and the switchboard-channel bridge
# + .mcp.json present in CHANNEL_WORKDIR.
#
# AGENT-ID / collision (DEPLOY-05): SWITCHBOARD_CHANNEL_AGENT_ID below is the bus id
# this responder long-polls. It MUST be the sole consumer of that id's inbox — either a
# dedicated id, or the id the retired claude-code-agent daemon used (disable that unit first).

[Unit]
Description=Claude Code Switchboard Channel Responder (headless --channels)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=%h/switchboard-channel
Environment="HOME=%h"
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="TERM=xterm-256color"
Environment="CHANNEL_WORKDIR=%h/switchboard-channel"
# --- deployment config (Phase 11): responder takes over the Billy id; the cold
#     claude-code-agent daemon is stopped+disabled so Billy's inbox has one consumer. ---
Environment="SWITCHBOARD_CHANNEL_AGENT_ID=Billy"
Environment="SWITCHBOARD_ALLOWED_SENDERS=Claude,Fred,Hermes,node24-upgrade,claude-code"
# Phase 12: bound context hourly via in-session /compact (keeps a summary; no downtime).
# Fallback if /compact ever wedges: drop this to 0 and uncomment RuntimeMaxSec below for
# clean hourly session rotation (Restart=always respawns; ~seconds downtime).
Environment="CHANNEL_COMPACT_INTERVAL_SEC=3600"
# Phase 13 / M1: restrict tool surface (untrusted bus content can't run shell or read secrets).
# Set to empty string only if you fully trust every allowlisted sender and want Billy to actuate.
Environment="CHANNEL_DISALLOWED_TOOLS=Bash Edit Write MultiEdit NotebookEdit"
ExecStart=/usr/bin/python3 %h/switchboard-channel/run-channel-session.py
Restart=always
RestartSec=10
# RuntimeMaxSec=3900   # rotation fallback (slightly >1h so /compact runs first)
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
